# unit/lavalink.service
[Unit]
Description=Tempo Bot, lavalink service config...
After=network-online.target

[Service]
User=albie
TimeoutStartSec=0
Type=simple
KillMode=process
WorkingDirectory=/home/albie/tempo/prod/
ExecStart=/bin/java -jar Lavalink.jar
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
